# spectratic-trading.service
[Unit]
Description=Spectratic Virtual Trading Engine
After=network.target spectratic-pipeline.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/us-stock-analyzer
ExecStart=/home/ubuntu/us-stock-analyzer/venv/bin/python -m trading.engine
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
