# spectratic-pipeline.service
[Unit]
Description=Spectratic Realtime Data Pipeline
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/us-stock-analyzer
ExecStart=/home/ubuntu/us-stock-analyzer/venv/bin/python -m data.pipeline
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
